# power_switch/files/switch
#!/bin/sh /etc/rc.common
#Automatic startup button power program
START = 65

start(){
    /usr/bin/switch 
}
